# install/sfn/safe-networking.conf
[program:safe-networking]
command=/home/pan/safe-networking-sp/sfn >/home/pan/safe-networking-sp/log/console-`date '+%Y-%m-%d-%H:%M:%S'`.log 2>&1
autostart=true
autorestart=true
